# home/dot_config/systemd/user/claude-code-router.service.tmpl
[Unit]
Description=Claude Code Router (tiered model routing proxy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Secrets file written by the chezmoi bootstrap script from
# ~/.claude/.credentials.json. Optional dash means "skip if missing".
EnvironmentFile=-%h/.claude-code-router/secrets.env
# mise shims live under ~/.local/share/mise/shims; PATH must include it.
Environment=PATH=%h/.local/share/mise/shims:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
{{- if get . "logPrompt" | default false }}
# Toggled via the logPrompt chezmoi flag — chezmoi apply re-renders this
# file and the run_onchange ccr-service script restarts the unit.
Environment=CCR_LOG_PROMPT=1
{{- end }}
{{- if eq (get . "cloudAuth" | default "apikey") "ccs-pro" }}
ExecStartPre=%h/.local/share/mise/shims/ccs cliproxy start
{{- end }}
ExecStart=%h/.local/share/mise/shims/ccr start
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
